# bin/initd
#!/bin/sh
### BEGIN INIT INFO
# Provides:       wolfpack
# Required-Start: $network $remote_fs $syslog
# Required-Stop:  $network $remote_fs $syslog
# Default-Start:  2 3 4 5
# Default-Stop:   1
# Description:    Wolfpack iptables configuration
### END INIT INFO

set -e # Exit on error status

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# DAEMON=/usr/local/bin/prax
# NAME=prax
DESC="Wolfpack iptables configuration"
HTTP_PORT=30666
HTTPS_PORT=30667

remove_rules_for_chain() {
  RULE_NUMS=`iptables -t nat -L $1 --line-numbers -n | egrep "($HTTP_PORT|$HTTPS_PORT)" | awk '{a[i++]=$0} END { for (j=i-1; j>=0;) print a[j--] }' | awk -F " " '{ print $1 }'`
  for RULE_NUM in $RULE_NUMS ; do
    iptables -t nat -D $1 $RULE_NUM
  done
}

case "$1" in
  start)
    if [ `iptables -t nat -L -n | egrep "($HTTP_PORT|$HTTPS_PORT)" | wc -l` -eq 0 ] ; then
      # for remote connections
      iptables -t nat -I PREROUTING -p tcp --dport 80 -j REDIRECT --to-ports $HTTP_PORT
      iptables -t nat -I PREROUTING -p tcp --dport 443 -j REDIRECT --to-ports $HTTPS_PORT
      # for local connections
      iptables -t nat -A OUTPUT -p tcp -o lo --dport 80  -j REDIRECT --to-ports $HTTP_PORT
      iptables -t nat -A OUTPUT -p tcp -o lo --dport 443 -j REDIRECT --to-ports $HTTPS_PORT
    fi
    ;;
  stop)
    remove_rules_for_chain 'OUTPUT'
    remove_rules_for_chain 'PREROUTING'
    ;;
  restart)
    $0 stop
    $0 start
    ;;
esac
